# systemd/tg_schedule_bot.service
[Unit]
Description=Telegram Schedule Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=schedulebot
Group=schedulebot
WorkingDirectory=/opt/tg_schedule_bot
EnvironmentFile=/opt/tg_schedule_bot/.env
ExecStart=/opt/tg_schedule_bot/.venv/bin/python /opt/tg_schedule_bot/main.py
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/tg_schedule_bot/data /opt/tg_schedule_bot/logs

[Install]
WantedBy=multi-user.target
