# systemd/tg_schedule_bot.service
[Unit]
Description=Telegram adapter for ops-core
After=network-online.target ops_core.service
Wants=network-online.target
Requires=ops_core.service

[Service]
Type=simple
User=schedulebot
Group=schedulebot
WorkingDirectory=/opt/tg_schedule_bot
EnvironmentFile=/opt/tg_schedule_bot/.env
ExecStart=/opt/tg_schedule_bot/.venv/bin/python /opt/tg_schedule_bot/main.py
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
